# elife/config/etc-init-uwsgi-loris.conf
description "uWSGI server for Loris"
start on (local-filesystems and net-device-up IFACE!=lo)
respawn
exec /opt/loris/venv/bin/uwsgi --ini /etc/loris2/uwsgi.ini

